# deploy/systemd/pulseops-api.service
[Unit]
Description=PulseOps registration API
After=network.target

[Service]
Type=simple
User=sasha
WorkingDirectory=/srv/www/pulseops.stage.dev/api
Environment=PORT=8106
Environment=PRODUCT_NAME=PulseOps
EnvironmentFile=/etc/skyyware/registration-mail.env
ExecStart=/usr/bin/node /srv/www/pulseops.stage.dev/api/registration-server.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
